# stock-ai/deployment/systemd/ai-trading-ai.service
[Unit]
Description=AI Trading System - Ai
After=network.target

[Service]
Type=simple
User=trading
Group=trading
WorkingDirectory=/opt/ai-trading
Environment=PYTHONPATH=/opt/ai-trading
ExecStart=python src/ai/ai_service.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
